# supervisor-tsdb.conf
[program:opentsdb]
command=/opt/sei-bin/start_opentsdb.sh
stopwaitsecs=60
