# supervisor/hbase.conf
[program:hbase]
command=/opt/sei-bin/start_hbase.sh
stopwaitsecs=60
